# app/supervisord.conf
[supervisord] 
nodaemon=true

[program:smbd]
command=/usr/sbin/smbd -FS

